# dev.nymann.hidecursor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTD/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.nymann.hidecursor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/hidecursor.app/Contents/MacOS/hidecursor</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/hidecursor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/hidecursor.err</string>
</dict>
</plist>
